# services/pull_data.timer
[Unit]
Description = Fetches data from the backend every week

[Timer]
OnCalendar = Weekly
Persistent = True

[Install]
WantedBy = timers.target
